# provisioning/Ansible/roles/application-deploy/templates/supervisor.conf
[program:demo-a-b]
command=/usr/bin/java -jar {{ app.basedir }}/demo-a-b/demo-a-b.jar
user={{ app.user }}
autostart=true
autorestart=true
startsecs=10
startretries=3
